# ansible/roles/lk-review-ingester/files/init.lk-review-ingester.conf
start on starting launchkit
stop on stopped launchkit

respawn
respawn limit 10 90

script
  cd /vagrant
  python review_ingester.py
end script