# ansible/roles/lk-django/files/init.lk-django.conf
start on (starting launchkit and started postgres)
stop on stopped launchkit

respawn
respawn limit 10 90

script
  cd /vagrant
  python manage.py runserver 0.0.0.0:9101
end script